# root/etc/s6-overlay/s6-rc.d/svc-sealskin/run
#!/usr/bin/with-contenv bash

# Run SealSkin
exec s6-setuidgid abc \
  bash -c 'cd /opt/sealskin/server && python3 main.py'
